# ansible-setup/roles/create-exercises/files/buchstabe@.service
[Unit]
Description=Service für Buchstabe

[Service]
Type=oneshot
ExecStart=/bin/bash -c "echo d > /home/%i/buchstabe"
User=%i
